# nginx/supervisord/conf.supervisord/cron.conf
[program:cron]
command=/usr/sbin/crond -f -l 2
stdout_syslog=true
stderr_syslog=true
autorestart=false
startentries=0
